# mowbot.service
[Unit]
Description=MowBot GUI + Stack
After=docker.service
Requires=docker.service

[Service]
Type=simple
# User/Group are overwritten by scripts/install.sh to match the installing account.
User=farmbot
Group=farmbot
# User= does not inherit docker group membership unless listed here (docker.sock access).
SupplementaryGroups=docker
WorkingDirectory=%h/mowbot-setup
Environment=HOME=%h
ExecStart=/usr/bin/docker compose --env-file mowbot.env -f docker-compose.yml up mowbot_gui
ExecStop=/usr/bin/docker compose --env-file mowbot.env -f docker-compose.yml stop -t 30 mowbot_gui
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
